# resources/docker/crontab
1,11,21,31,41,51 * * * * app /app/resources/docker/etl.sh
