# resources/docker/crontab
10 * * * * app /app/resources/docker/etl.shecho "Hello world" >> /var/log/cron.log 2>&1
